# deploy/s6-overlay/s6-rc.d/postgres/run
#!/command/execlineb -P
with-contenv
foreground { s6-echo "postgres: starting" }
s6-setuidgid postgres
/usr/lib/postgresql/16/bin/postgres
  -D /data/postgres
  -c config_file=/data/postgres/postgresql.conf
