# deploy/s6-overlay/s6-rc.d/redis/run
#!/command/execlineb -P
with-contenv
foreground { s6-echo "redis: starting" }
foreground { install -d -o sda -g sda -m 0750 /data/redis }
s6-setuidgid sda
/usr/bin/redis-server /etc/sda/redis.conf
