# deploy/s6-overlay/s6-rc.d/app/run
#!/command/execlineb -P
with-contenv
s6-setuidgid sda
/opt/sda/bin/app
