# deploy/s6-overlay/s6-rc.d/nats/run
#!/command/execlineb -P
with-contenv
foreground { install -d -o sda -g sda -m 0750 /data/nats }
s6-setuidgid sda
/usr/local/bin/nats-server -c /etc/sda/nats.conf
